# launchd/com.chrisobrien-ai.pr-watch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  pr-watch LaunchAgent — polls GitHub every 5 min for PRs needing an auto-fix.
  TEMPLATE: install.sh materializes this to ~/Library/LaunchAgents/ with __HOME__
  replaced by the real path (a symlinked plist can fail to bootstrap). The agent is
  INERT until `pr-watch enable` (poll() no-ops without ~/.config/pr-watch/enabled),
  so installing it does not, by itself, start auto-fixing anything.

  PATH is set explicitly: launchd hands an agent a minimal PATH, but pr-watch shells
  out to gh / git / tmux, which live in Homebrew. The pane shells it spawns re-source
  .zshrc, but the watcher process itself needs Homebrew on PATH to find gh.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.chrisobrien-ai.pr-watch</string>

    <key>ProgramArguments</key>
    <array>
        <string>__HOME__/bin/pr-watch</string>
        <string>poll</string>
    </array>

    <key>StartInterval</key>
    <integer>300</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/pr-watch.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/pr-watch.log</string>
</dict>
</plist>
